# packaging/systemd/resguardd.service
[Unit]
Description=Resguard freeze watchdog daemon
Documentation=https://github.com/resguard/resguard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/resguardd --config /etc/resguard/resguardd.yml
Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
